# config/includes.chroot/etc/systemd/system/permission-hardening.service
[Unit]
Description=Permission System Hardening Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/permissionhardening
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 